# systemd/user/taey-display-2.service
[Unit]
Description=Taey Display :2 (ChatGPT)
After=default.target

[Service]
Type=simple
Environment=DISPLAY=:2
ExecStartPre=/usr/bin/bash -lc 'rm -f /tmp/.X2-lock /tmp/.X11-unix/X2 /tmp/a11y_bus_:2 /tmp/firefox_pid_:2; pkill -f "Xvfb :2 " 2>/dev/null || true; pkill -f "/tmp/ff-profile-mira-chatgpt" 2>/dev/null || true; /usr/bin/Xvfb :2 -screen 0 1920x1080x24 -ac -noreset >/tmp/taey-display-2-xvfb.log 2>&1 &'
ExecStart=/usr/bin/dbus-run-session -- /usr/bin/bash -lc '/usr/libexec/at-spi-bus-launcher --launch-immediately & sleep 1; A11Y_ADDR=$(xprop -display :2 -root AT_SPI_BUS 2>/dev/null | sed "s/.*= \\"//;s/\\"$//"); echo "${A11Y_ADDR}" > /tmp/a11y_bus_:2; openbox & sleep 1; mkdir -p /tmp/ff-profile-mira-chatgpt; GTK_USE_PORTAL=0 LIBGL_ALWAYS_SOFTWARE=1 firefox --no-remote --profile /tmp/ff-profile-mira-chatgpt "https://chatgpt.com" & FIREFOX_PID=$!; echo "${FIREFOX_PID}" > /tmp/firefox_pid_:2; wait ${FIREFOX_PID}'
ExecStopPost=/usr/bin/bash -lc 'pkill -f "Xvfb :2 " 2>/dev/null || true; rm -f /tmp/.X2-lock /tmp/.X11-unix/X2 /tmp/a11y_bus_:2 /tmp/firefox_pid_:2'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
