# systemd/user/taey-display-5.service
[Unit]
Description=Taey Display :5 (Grok)
After=taey-xvfb@5.service
Requires=taey-xvfb@5.service

[Service]
Type=simple
Environment=DISPLAY=:5
ExecStart=/usr/bin/dbus-run-session -- /usr/bin/bash -lc 'echo "$DBUS_SESSION_BUS_ADDRESS" > /tmp/dbus_session_bus_:5; /usr/libexec/at-spi-bus-launcher --launch-immediately & sleep 2; A11Y_ADDR=$(xprop -display :5 -root AT_SPI_BUS 2>/dev/null | sed "s/.*= \"//;s/\"$//"); echo "${A11Y_ADDR}" > /tmp/a11y_bus_:5; openbox & sleep 1; /usr/bin/x11vnc -display :5 -rfbport 5905 -rfbauth $HOME/.taey/vnc_passwd -forever -shared -bg -o /tmp/x11vnc-5.log; mkdir -p $HOME/.taey/profiles/grok; GTK_USE_PORTAL=0 LIBGL_ALWAYS_SOFTWARE=1 /usr/lib/firefox/firefox --no-remote --profile $HOME/.taey/profiles/grok "https://grok.com" & FIREFOX_PID=$!; echo "${FIREFOX_PID}" > /tmp/firefox_pid_:5; wait ${FIREFOX_PID}'
ExecStopPost=/usr/bin/bash -lc 'pkill -f "x11vnc.*rfbport 5905" 2>/dev/null || true; rm -f /tmp/a11y_bus_:5 /tmp/dbus_session_bus_:5 /tmp/firefox_pid_:5'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
