# systemd/user/taey-display-3.service
[Unit]
Description=Taey Display :3 (Claude)
After=network.target

[Service]
Type=simple
Environment="DISPLAY=:3"
EnvironmentFile=/home/mira/taeys-hands/.env
Environment="PROFILE=/tmp/ff-profile-%H-claude-3"

ExecStartPre=/bin/sh -c 'rm -f /tmp/a11y_bus_:3 && mkdir -p ${PROFILE}'

# Start blocking isolated D-Bus. Write address to standard file.
ExecStart=/usr/bin/dbus-run-session -- sh -c '\
    echo "export DBUS_SESSION_BUS_ADDRESS=$DBUS_SESSION_BUS_ADDRESS" > /tmp/a11y_bus_:3 && \
    /usr/libexec/at-spi-bus-launcher & \
    /usr/libexec/at-spi2-registryd & \
    Xvfb :3 -screen 0 1920x1080x24 & \
    sleep 2 && \
    exec firefox --profile ${PROFILE} --new-instance "https://claude.ai/new?incognito"'

Restart=always
RestartSec=5

[Install]
WantedBy=default.target
