# systemd/user/taey-display-4.service
[Unit]
Description=Taey Display :4 (Gemini)
After=network.target

[Service]
Type=simple
Environment="DISPLAY=:4"
EnvironmentFile=/home/mira/taeys-hands/.env
Environment="PROFILE=/tmp/ff-profile-%H-gemini-4"

ExecStartPre=/bin/sh -c 'rm -f /tmp/a11y_bus_:4 && mkdir -p ${PROFILE}'

# Start blocking isolated D-Bus. Write address to standard file.
ExecStart=/usr/bin/dbus-run-session -- sh -c '\
    echo "export DBUS_SESSION_BUS_ADDRESS=$DBUS_SESSION_BUS_ADDRESS" > /tmp/a11y_bus_:4 && \
    /usr/libexec/at-spi-bus-launcher & \
    /usr/libexec/at-spi2-registryd & \
    Xvfb :4 -screen 0 1920x1080x24 & \
    sleep 2 && \
    exec firefox --profile ${PROFILE} --new-instance "https://gemini.google.com/app"'

Restart=always
RestartSec=5

[Install]
WantedBy=default.target
