# systemd/user/taey-display-6.service
[Unit]
Description=Taey Display :6 (Perplexity)
After=network.target

[Service]
Type=simple
Environment="DISPLAY=:6"
EnvironmentFile=/home/mira/taeys-hands/.env
Environment="PROFILE=/tmp/ff-profile-%H-perplexity-6"

ExecStartPre=/bin/sh -c 'rm -f /tmp/a11y_bus_:6 && mkdir -p ${PROFILE}'

# Start blocking isolated D-Bus. Write address to standard file.
ExecStart=/usr/bin/dbus-run-session -- sh -c '\
    echo "export DBUS_SESSION_BUS_ADDRESS=$DBUS_SESSION_BUS_ADDRESS" > /tmp/a11y_bus_:6 && \
    /usr/libexec/at-spi-bus-launcher & \
    /usr/libexec/at-spi2-registryd & \
    Xvfb :6 -screen 0 1920x1080x24 & \
    sleep 2 && \
    exec firefox --profile ${PROFILE} --new-instance "https://www.perplexity.ai/"'

Restart=always
RestartSec=5

[Install]
WantedBy=default.target
