# systemd/user/taey-display-3.service
[Unit]
Description=Taey Display :3 (Claude)
After=taey-xvfb@3.service
Requires=taey-xvfb@3.service

[Service]
Type=simple
Environment=DISPLAY=:3
Environment=TAEY_REPO=%h/taeys-hands
ExecStart=/usr/bin/dbus-run-session -- /usr/bin/bash -lc 'echo "$$DBUS_SESSION_BUS_ADDRESS" > /tmp/dbus_session_bus_:3; /usr/libexec/at-spi-bus-launcher --launch-immediately & A11Y_ADDR=""; for i in $$(seq 1 20); do TMP_ADDR=$$(xprop -display :3 -root AT_SPI_BUS 2>/dev/null | sed "s/.*= \"//;s/\"$$//"); case "$$TMP_ADDR" in unix:path=*|unix:abstract=*) A11Y_ADDR="$$TMP_ADDR"; echo "$$A11Y_ADDR" > /tmp/a11y_bus_:3 && break ;; esac; sleep 1; done; if [[ -z "$$A11Y_ADDR" ]]; then TMP_ADDR=$$(xprop -display :3 -root AT_SPI_BUS 2>/dev/null | sed "s/.*= \"//;s/\"$$//"); case "$$TMP_ADDR" in unix:path=*|unix:abstract=*) A11Y_ADDR="$$TMP_ADDR"; echo "$$A11Y_ADDR" > /tmp/a11y_bus_:3 ;; esac; fi; AT_SPI_BUS_ADDRESS="$${A11Y_ADDR}" /usr/libexec/at-spi2-registryd & sleep 1; openbox & sleep 1; /usr/bin/x11vnc -display :3 -rfbport 5903 -rfbauth $$HOME/.taey/vnc_passwd -forever -shared -bg -o /tmp/x11vnc-3.log; mkdir -p $$HOME/.taey/profiles/ff-profile-claude; cp "$${TAEY_REPO}/systemd/user/firefox-user.js" $$HOME/.taey/profiles/ff-profile-claude/user.js; echo "user_pref(\"accessibility.force_disabled\", -1);" >> $$HOME/.taey/profiles/ff-profile-claude/user.js; GTK_USE_PORTAL=0 LIBGL_ALWAYS_SOFTWARE=1 GNOME_ACCESSIBILITY=1 GTK_MODULES=gail:atk-bridge AT_SPI_BUS_ADDRESS="$${A11Y_ADDR}" /usr/lib/firefox/firefox --no-remote --profile $$HOME/.taey/profiles/ff-profile-claude "https://claude.ai/new" & FIREFOX_PID=$$!; if [[ -z "$$FIREFOX_PID" ]]; then for _ in $$(seq 1 10); do FIREFOX_PID=$$(pgrep -P $$$$ -o -x firefox 2>/dev/null || true); [[ -n "$$FIREFOX_PID" ]] && break; sleep 1; done; fi; ( L_ADDR="$$A11Y_ADDR"; while kill -0 "$${FIREFOX_PID}" >/dev/null 2>&1; do C_ADDR=$$(xprop -display :3 -root AT_SPI_BUS 2>/dev/null | sed "s/.*= \"//;s/\"$$//"); if [[ "$$C_ADDR" == unix:* ]]; then if [[ "$$C_ADDR" != "$$L_ADDR" || ! -s /tmp/a11y_bus_:3 ]]; then printf "%%s\n" "$$C_ADDR" > /tmp/a11y_bus_:3.tmp && mv /tmp/a11y_bus_:3.tmp /tmp/a11y_bus_:3; L_ADDR="$$C_ADDR"; fi; fi; sleep 2; done ) & [ -s /tmp/a11y_bus_:3 ] || echo "WARNING: AT-SPI bus capture failed for :3" >&2; echo "$${FIREFOX_PID}" > /tmp/firefox_pid_:3; if [[ -n "$${FIREFOX_PID}" ]]; then wait $${FIREFOX_PID}; else wait; fi'
ExecStopPost=/usr/bin/bash -lc 'pkill -f "x11vnc.*rfbport 5903" 2>/dev/null || true; rm -f /tmp/dbus_session_bus_:3 /tmp/firefox_pid_:3'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
