# systemd/user/taey-display-5.service
[Unit]
Description=Taey Display :5 (Grok)
After=network.target

[Service]
Type=simple
Environment="DISPLAY=:5"
EnvironmentFile=/home/mira/taeys-hands/.env
Environment="PROFILE=/tmp/ff-profile-%H-grok-5"

ExecStartPre=/bin/sh -c 'rm -f /tmp/a11y_bus_:5 && mkdir -p ${PROFILE}'

# Start blocking isolated D-Bus. Write address to standard file.
ExecStart=/usr/bin/dbus-run-session -- sh -c '\
    echo "export DBUS_SESSION_BUS_ADDRESS=$DBUS_SESSION_BUS_ADDRESS" > /tmp/a11y_bus_:5 && \
    /usr/libexec/at-spi-bus-launcher & \
    /usr/libexec/at-spi2-registryd & \
    Xvfb :5 -screen 0 1920x1080x24 & \
    sleep 2 && \
    exec firefox --profile ${PROFILE} --new-instance "https://grok.com/"'

Restart=always
RestartSec=5

[Install]
WantedBy=default.target
