# systemd/user/taey-display-2.service
[Unit]
Description=Taey Display :2 (ChatGPT)
After=taey-xvfb@2.service
Requires=taey-xvfb@2.service

[Service]
Type=simple
Environment=DISPLAY=:2
ExecStart=/usr/bin/dbus-run-session -- /usr/bin/bash -lc 'echo "$DBUS_SESSION_BUS_ADDRESS" > /tmp/dbus_session_bus_:2; /usr/libexec/at-spi-bus-launcher --launch-immediately & sleep 2; A11Y_ADDR=$(xprop -display :2 -root AT_SPI_BUS 2>/dev/null | sed "s/.*= \"//;s/\"$//"); echo "${A11Y_ADDR}" > /tmp/a11y_bus_:2; openbox & sleep 1; /usr/bin/x11vnc -display :2 -rfbport 5902 -rfbauth $HOME/.taey/vnc_passwd -forever -shared -bg -o /tmp/x11vnc-2.log; mkdir -p $HOME/.taey/profiles/chatgpt; GTK_USE_PORTAL=0 LIBGL_ALWAYS_SOFTWARE=1 /usr/lib/firefox/firefox --no-remote --profile $HOME/.taey/profiles/chatgpt "https://chatgpt.com" & FIREFOX_PID=$!; echo "${FIREFOX_PID}" > /tmp/firefox_pid_:2; wait ${FIREFOX_PID}'
ExecStopPost=/usr/bin/bash -lc 'pkill -f "x11vnc.*rfbport 5902" 2>/dev/null || true; rm -f /tmp/a11y_bus_:2 /tmp/dbus_session_bus_:2 /tmp/firefox_pid_:2'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
