# systemd/user/taey-display-2.service
[Unit]
Description=Taey Display :2 (ChatGPT)
After=network.target

[Service]
Type=simple
Environment="DISPLAY=:2"
EnvironmentFile=/home/mira/taeys-hands/.env
Environment="PROFILE=/tmp/ff-profile-%H-chatgpt-2"

ExecStartPre=/bin/sh -c 'rm -f /tmp/a11y_bus_:2 && mkdir -p ${PROFILE}'

# Start blocking isolated D-Bus. Write address to standard file.
ExecStart=/usr/bin/dbus-run-session -- sh -c '\
    echo "export DBUS_SESSION_BUS_ADDRESS=$DBUS_SESSION_BUS_ADDRESS" > /tmp/a11y_bus_:2 && \
    /usr/libexec/at-spi-bus-launcher & \
    /usr/libexec/at-spi2-registryd & \
    Xvfb :2 -screen 0 1920x1080x24 & \
    sleep 2 && \
    exec firefox --profile ${PROFILE} --new-instance "https://chatgpt.com/?temporary-chat=true"'

Restart=always
RestartSec=5

[Install]
WantedBy=default.target
