# systemd/user/taey-display-3.service
[Unit]
Description=Taey Display :3 (Claude)
After=taey-xvfb@3.service
Requires=taey-xvfb@3.service

[Service]
Type=simple
Environment=DISPLAY=:3
ExecStart=/usr/bin/dbus-run-session -- /usr/bin/bash -lc 'echo "$DBUS_SESSION_BUS_ADDRESS" > /tmp/dbus_session_bus_:3; /usr/libexec/at-spi-bus-launcher --launch-immediately & sleep 2; A11Y_ADDR=$(xprop -display :3 -root AT_SPI_BUS 2>/dev/null | sed "s/.*= \"//;s/\"$//"); echo "${A11Y_ADDR}" > /tmp/a11y_bus_:3; openbox & sleep 1; /usr/bin/x11vnc -display :3 -rfbport 5903 -rfbauth $HOME/.taey/vnc_passwd -forever -shared -bg -o /tmp/x11vnc-3.log; mkdir -p $HOME/.taey/profiles/claude; GTK_USE_PORTAL=0 LIBGL_ALWAYS_SOFTWARE=1 /usr/lib/firefox/firefox --no-remote --profile $HOME/.taey/profiles/claude "https://claude.ai/new" & FIREFOX_PID=$!; echo "${FIREFOX_PID}" > /tmp/firefox_pid_:3; wait ${FIREFOX_PID}'
ExecStopPost=/usr/bin/bash -lc 'pkill -f "x11vnc.*rfbport 5903" 2>/dev/null || true; rm -f /tmp/a11y_bus_:3 /tmp/dbus_session_bus_:3 /tmp/firefox_pid_:3'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
