# systemd/user/taey-bus-watcher@.service
[Unit]
Description=Taey AT-SPI bus pointer watcher for display :%i
Requires=taey-display-%i.service
After=taey-display-%i.service
PartOf=taey-display-%i.service

[Service]
Type=simple
EnvironmentFile=%h/.taey/machine.env
Environment=DISPLAY=:%i
ExecStart=/usr/bin/bash -lc 'set -Eeuo pipefail; : "$${TAEY_REPO:?TAEY_REPO}"; exec "$${TAEY_REPO}/scripts/bus_watcher.sh" %i'
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
