# systemd/user/taey-display-6.service
[Unit]
Description=Taey Display :6 (Perplexity)
After=taey-xvfb@6.service
Requires=taey-xvfb@6.service

[Service]
Type=simple
Environment=DISPLAY=:6
ExecStart=/usr/bin/dbus-run-session -- /usr/bin/bash -lc 'echo "$DBUS_SESSION_BUS_ADDRESS" > /tmp/dbus_session_bus_:6; /usr/libexec/at-spi-bus-launcher --launch-immediately & sleep 2; A11Y_ADDR=$(xprop -display :6 -root AT_SPI_BUS 2>/dev/null | sed "s/.*= \"//;s/\"$//"); echo "${A11Y_ADDR}" > /tmp/a11y_bus_:6; openbox & sleep 1; /usr/bin/x11vnc -display :6 -rfbport 5906 -rfbauth $HOME/.taey/vnc_passwd -forever -shared -bg -o /tmp/x11vnc-6.log; mkdir -p $HOME/.taey/profiles/perplexity; cp $HOME/taeys-hands/systemd/user/firefox-user.js $HOME/.taey/profiles/perplexity/user.js; GTK_USE_PORTAL=0 LIBGL_ALWAYS_SOFTWARE=1 /usr/lib/firefox/firefox --no-remote --profile $HOME/.taey/profiles/perplexity "https://www.perplexity.ai/" & FIREFOX_PID=$!; echo "${FIREFOX_PID}" > /tmp/firefox_pid_:6; wait ${FIREFOX_PID}'
ExecStopPost=/usr/bin/bash -lc 'pkill -f "x11vnc.*rfbport 5906" 2>/dev/null || true; rm -f /tmp/a11y_bus_:6 /tmp/dbus_session_bus_:6 /tmp/firefox_pid_:6'
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
