# scripts/conductor-dashboard.service
[Unit]
Description=The Conductor's Console — AI Orchestration Dashboard
After=network.target redis.service neo4j.service

[Service]
Type=simple
User=spark
WorkingDirectory=/home/spark/taeys-hands
ExecStart=/home/spark/.local/bin/uvicorn dashboard.app:app --host 0.0.0.0 --port 5001
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
