# orchestration/heartbeats.service
[Unit]
Description=Agent Heartbeat Daemons — Local Agents
After=network.target conductor.service
Requires=conductor.service

[Service]
Type=simple
User=spark
WorkingDirectory=/home/spark/taeys-hands
ExecStart=/bin/bash -c '\
  python3 orchestration/agent_beat.py claude-weaver "idle" & \
  python3 orchestration/agent_beat.py codex-cli "idle" & \
  python3 orchestration/agent_beat.py gemini-cli "idle" & \
  wait'
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
